# scripts/com.iris.claude-code-vault-gardener.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.iris.claude-code-vault-gardener</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>-lc</string>
		<string>/Volumes/AI_Workspace/iris_studio/scripts/run_claude_job.sh vault-gardener /Volumes/AI_Workspace/iris_studio/routines/vault-gardener.prompt</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>StandardErrorPath</key>
	<string>/Users/steve/iris_studio/logs/claude-code-vault-gardener.stderr.log</string>
	<key>StandardOutPath</key>
	<string>/Users/steve/iris_studio/logs/claude-code-vault-gardener.stdout.log</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>4</integer>
		<key>Minute</key>
		<integer>10</integer>
	</dict>
	<!-- The vault gardener is file-descriptor-heavy (vault-wide 24h diff +
	     per-note relatedness search + filesystem/obsidian MCP servers). Under
	     launchd's default soft RLIMIT_NOFILE of 256 the real 2026-06-16 04:10
	     fire died with "low max file descriptors (Current limit: 256)". Raise
	     the soft limit so the scheduled fire has headroom (the manual kickstart
	     succeeded only because it inherited an interactive shell's larger limit). -->
	<key>SoftResourceLimits</key>
	<dict>
		<key>NumberOfFiles</key>
		<integer>8192</integer>
	</dict>
</dict>
</plist>
