# scripts/com.iris.claude-code-caption-sweep.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>JOB_QUIET_OK</key>
		<string>1</string>
	</dict>
	<key>Label</key>
	<string>com.iris.claude-code-caption-sweep</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Volumes/AI_Workspace/iris_studio/scripts/run_job.sh</string>
		<string>caption-sweep</string>
		<string>/Volumes/AI_Workspace/iris_studio/.venv/bin/python3</string>
		<string>/Volumes/AI_Workspace/iris_studio/scripts/sweep_captions.py</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>SoftResourceLimits</key>
	<dict>
		<key>NumberOfFiles</key>
		<integer>8192</integer>
	</dict>
	<key>StandardErrorPath</key>
	<string>/Users/steve/iris_studio/logs/claude-code-caption-sweep.stderr.log</string>
	<key>StandardOutPath</key>
	<string>/Users/steve/iris_studio/logs/claude-code-caption-sweep.stdout.log</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>4</integer>
		<key>Minute</key>
		<integer>20</integer>
	</dict>
	<key>WorkingDirectory</key>
	<string>/Volumes/AI_Workspace/iris_studio</string>
</dict>
</plist>
