# scripts/com.iris.claude-code-drive-gardener.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.iris.claude-code-drive-gardener</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>-lc</string>
		<string>cd /Users/steve/Documents/3SK/outputs &amp;&amp; /opt/homebrew/bin/claude --print --dangerously-skip-permissions "$(cat /Volumes/AI_Workspace/iris_studio/routines/drive-gardener.prompt)" 2&gt;&amp;1 | tee -a /Users/steve/iris_studio/logs/claude-code-drive-gardener.log</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>StandardErrorPath</key>
	<string>/Users/steve/iris_studio/logs/claude-code-drive-gardener.stderr.log</string>
	<key>StandardOutPath</key>
	<string>/Users/steve/iris_studio/logs/claude-code-drive-gardener.stdout.log</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>3</integer>
		<key>Minute</key>
		<integer>0</integer>
		<key>Weekday</key>
		<integer>0</integer>
	</dict>
</dict>
</plist>
