# scripts/com.iris.claude-code-analytics-pull.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Build 4 weekly analytics FEEDBACK LOOP — installed 2026-06-19 (V1+V2 live).

  CLOSED LOOP, not just a data pull: this runs the `analytics-feedback` Claude
  routine, which (1) runs analytics_pull.py to fetch our channel's real metrics
  from the YouTube Analytics API, then (2) — ONLY when the feed has real data
  (status: ok) — dispatches the channel-analyst agent to produce the routed
  fixes for the content agents. A no-data week (pre-launch, or data still lagging
  ~48-72h after publish) correctly skips the analyst and reports a clean no-op;
  the analyst never runs on empty data and never fabricates.

  Runs via run_claude_job.sh (NOT run_job.sh): it's a `claude --print` routine, so
  the SAME wrapper that every other claude-code job uses handles logging, the
  ROUTINE_COMPLETE/INCOMPLETE sentinel → Telegram alert, and the retry queue.

  Schedule: Monday 5:30 AM ET — diagnosis fresh before the weekly production cycle.

  Manage:
    launchctl bootout gui/$(id -u)/com.iris.claude-code-analytics-pull 2>/dev/null
    cp scripts/com.iris.claude-code-analytics-pull.plist ~/Library/LaunchAgents/
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.iris.claude-code-analytics-pull.plist
    launchctl kickstart -k gui/$(id -u)/com.iris.claude-code-analytics-pull   # test fire
-->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.iris.claude-code-analytics-pull</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>-lc</string>
		<string>/Volumes/AI_Workspace/iris_studio/scripts/run_claude_job.sh analytics-pull /Volumes/AI_Workspace/iris_studio/routines/analytics-feedback.prompt</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>StandardErrorPath</key>
	<string>/Users/steve/iris_studio/logs/analytics-pull.stderr.log</string>
	<key>StandardOutPath</key>
	<string>/Users/steve/iris_studio/logs/analytics-pull.stdout.log</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>5</integer>
		<key>Minute</key>
		<integer>30</integer>
		<key>Weekday</key>
		<integer>1</integer>
	</dict>
</dict>
</plist>
